# aports/device/device-mozilla-flame/50-firmware.rules
SUBSYSTEM=="firmware", ACTION=="add", RUN+="/usr/bin/firmware.sh"